# packages/genx-ai-camera-panel/genx_ai_camera_panel-2.2.0-cp312-cp312-manylinux_2_17_aarch64.manylinux2014_aarch64.whl/genx_ai_camera_panel/service/genx-ai-camera-panel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>genx-ai-camera-panel</string>
        <key>ProgramArguments</key>
        <array>
            <string>python3</string>
            <string>-m</string>
            <string>genx_ai_camera_panel.main</string>
            <string>listen</string>
        </array>
        <key>RunAtLoad</key>
        <true/>
        <key>KeepAlive</key>
        <true/>
        <key>StandardErrorPath</key>
        <string>/tmp/genx-ai-camera-panel.err</string>
        <key>StandardOutPath</key>
        <string>/tmp/genx-ai-camera-panel.out</string>
    </dict>
</plist>
