# packages/genx-ai-camera-panel/genx_ai_camera_panel-2.2.0-cp312-cp312-manylinux_2_17_aarch64.manylinux2014_aarch64.whl/genx_ai_camera_panel/service/genx-ai-camera-panel.service
[Unit]
Description=GenX AI Camera Panel Service
After=network.target

[Service]
ExecStart=/opt/venv/bin/genx-camera-panel listen
Restart=always
Type=simple

[Install]
WantedBy=multi-user.target
