# packages/rpi-metrics-server/rpi_metrics_server-1.0.0.tar.gz/rpi_metrics_server-1.0.0/rpi_metrics/server/rpi-metricsd.service
[Unit]
Description=RPi Metrics Service Daemon
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/usr/share/rpi-metrics
ExecStart=/usr/share/rpi-metrics/Server/venv/bin/python /usr/share/rpi-metrics/Server/rpi_metrics_server.py
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target